# deploy/systemd/pao.service
[Unit]
Description=pao.brasume.com — calculadora de pão
After=network.target

[Service]
Type=simple
User=pao
Group=pao
WorkingDirectory=/opt/pao/app
EnvironmentFile=/opt/pao/app/.env
ExecStart=/usr/bin/node --env-file=.env dist-server/index.js
Restart=on-failure
RestartSec=3

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/pao
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
